# vagrant/files/etc/update-motd.d/50-esqulino
#!/bin/bash

printf "\nDefault Ports\n"
printf " * Local HTTP server : 9292\n"
printf " * Host  HTTP mapping: 9292\n"

echo ""
systemctl status esqulino
